# build/debian/prerm
#!/bin/sh

sudo systemctl stop hacienda
sudo systemctl disable hacienda
sudo systemctl daemon-reload
